# deploy/zippy.service
# Zippy systemd USER service (Benedict — rootless, no sudo needed).
#
#   mkdir -p ~/.config/systemd/user
#   cp deploy/zippy.service ~/.config/systemd/user/zippy.service
#   loginctl enable-linger                 # survive logout/reboot (self; no sudo)
#   systemctl --user daemon-reload
#   systemctl --user enable --now zippy
#
# Serves the Next.js production build on 127.0.0.1:3000; Tailscale Funnel
# (port 8443) exposes it publicly. `next start` forces NODE_ENV=production,
# which turns on Secure cookies.
[Unit]
Description=Zippy — EVE Online wormhole route planner (Next.js)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/paul/zippy
EnvironmentFile=/home/paul/zippy/.env.local
ExecStart=/usr/local/bin/pnpm start -H 127.0.0.1 -p 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
